# sources/supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
priority=1
command=/usr/local/bin/ovpn_run
autorestart=true

[program:nginx]
priority=2
command=/usr/sbin/nginx -g "daemon off;"
autorestart=true